# etc/supervisor/conf.d/mariadb.conf
[program:mariadb]
command=/usr/bin/mysqld --user=mysql --console

stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
